# host/etc/rc.local
#!/bin/bash
/opt/cyberpot/bin/updateip.sh
exit 0
